# docker/cron/crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

* * * * * docker exec itcapital-app php artisan schedule:run >> /var/log/cron.log 2>&1
